# systemd/controller-legend.service
[Unit]
Description=Controller Legend HUD
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/projects/-AI-controller./scripts/controller-legend.py
Restart=on-failure
RestartSec=3
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
