# systemd/ptt-pynput.service
[Unit]
Description=Push-to-Talk Dictation Listener (pynput)
After=graphical-session.target voice-bridge.service
Wants=graphical-session.target voice-bridge.service

[Service]
Type=simple
ExecStart=%h/ai-controller/.venv/bin/python3 %h/ai-controller/scripts/ptt_pynput.py
WorkingDirectory=%h/ai-controller/scripts
Restart=on-failure
RestartSec=3
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=XAUTHORITY=%h/.Xauthority

[Install]
WantedBy=default.target
