# systemd/voice-bridge.service
[Unit]
Description=Voice Bridge — PTT → Groq Whisper → CLAF → spd-say
After=claf.service
Wants=claf.service

[Service]
ExecStart=/usr/bin/python3 %h/projects/-AI-controller./scripts/voice_bridge.py
WorkingDirectory=%h/projects/-AI-controller.
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
